# app/packaging/rnf-backend/debian/rnf-backend.postinst
#!/bin/bash

set -e

/opt/rnf/scripts/db_bootstrap.sh

systemctl daemon-reload
systemctl enable rnf.service

echo "Starting RNF service"
systemctl start rnf.service